# deploy/auto-pull.service
[Unit]
Description=Storybook auto-pull from GitHub

[Service]
Type=oneshot
User=prometheus
ExecStart=/ssd/storybook/deploy/auto-pull.sh
